# roles/immutant/templates/immutant-upstart.conf.j2
description "This is an upstart job file for Immutant"
 
# Ubuntu doesn't yet emit the network-services event:
# https://bugs.launchpad.net/ubuntu/+source/upstart/+bug/701576
#
#start on started network-services
#stop on stopping network-services
 
start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]
# respawn
 
limit nofile 4096 4096

kill signal SIGINT

expect fork

script
bash << "EOF"
  su - immutant
  export IMMUTANT_HOME=/var/lib/immutant/current
  exec lein immutant run --clustered -b {{ ec2_private_ip_address }} -Djboss.port={{ immutant_port }} >> /var/log/immutant/immutant.log 2>&1
EOF
end script